# systemd/anomaly-consumer.service
# ⚠️ 경로 수정 필요: /home/ubuntu/finance_coin → 본인 환경에 맞게 변경

[Unit]
Description=Anomaly Detection Consumer (이상탐지 → TimescaleDB)
After=network.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/finance_coin
ExecStart=/home/ubuntu/finance_coin/.venv/bin/python consumers/anomaly_consumer.py
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5
EnvironmentFile=/home/ubuntu/finance_coin/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anomaly-consumer

[Install]
WantedBy=multi-user.target