# systemd/okx-producer.service
[Unit]
Description=OKX WebSocket Producer
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=300
StartLimitBurst=30

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/finance_coin
EnvironmentFile=/home/ubuntu/finance_coin/.env
Environment=PYTHONUNBUFFERED=1
TimeoutStartSec=0

ExecStartPre=/bin/bash -c 'until timeout 10 docker exec airflow-kafka-1 kafka-broker-api-versions --bootstrap-server localhost:9092 >/dev/null 2>&1; do echo "[okx-producer] waiting for kafka"; sleep 5; done'
ExecStart=/home/ubuntu/finance_coin/.venv/bin/python producers/okx_producer.py

Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=okx-producer

[Install]
WantedBy=multi-user.target
